# cmd/maintenance/owui-maintenance@.service
[Unit]
Description=Open WebUI maintenance task: %i
After=open-webui.service
Before=bc250-night-shutdown.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/bc250-llm-server/maintenance.env
ExecStart=/usr/libexec/bc250-llm-server/%i.sh
Nice=10
IOSchedulingClass=idle
