# maintenance/owui-backup-users.service
[Unit]
Description=Open WebUI identity-only backup
After=open-webui.service
Before=bc250-night-shutdown.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/bc250-llm-server/maintenance.env
ExecStart=/usr/libexec/bc250-llm-server/backup-users.sh
Nice=10
IOSchedulingClass=idle
