# maintenance/owui-prune.timer
[Unit]
Description=Daily Open WebUI upload prune

[Timer]
OnCalendar=*-*-* 18:10:00
Persistent=false
AccuracySec=1min

[Install]
WantedBy=timers.target
