# examples/raspi-wol/bc250-night-shutdown.service
# Target on BC-250: /etc/systemd/system/bc250-night-shutdown.service
[Unit]
Description=Suspend the LLM server only when idle
After=owui-maintenance@prune-uploads.service
After=owui-maintenance@backup-users.service owui-maintenance@backup-config.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/bc250-llm-server/maintenance.env
ExecStart=/usr/libexec/bc250-llm-server/safe-suspend.sh
